# raspi/ansible/roles/clinicalsupervisor/templates/clinicalsupervisor-supervisor.conf
[program:clinicalsupervisor]
command=/bin/bash -c '{{ ansible_distribution|get_clinicalsupervisor_dir }}/clinicalsupervisor/venv/bin/uwsgi --ini /etc/clinicalsupervisor/uwsgi.ini'

user={{ supervisor_user }}
numprocs=1
stdout_logfile=/var/log/clinicalsupervisor/clinicalsupervisor_stdout.log
stderr_logfile=/var/log/clinicalsupervisor/clinicalsupervisor_stderr.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
stopsignal=INT
priority=998
